# src/modules/hotspot_cb1/filesystem/root/usr/lib/systemd/system/create_ap.service
[Unit]
Description=Create AP Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/start_hotspot.sh
KillSignal=SIGINT
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target